# automagic-fan.service
[Unit]
Description=Automagic fan control

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /usr/local/bin/automagic-fan/fanctl.py
ExecStopPost=/bin/sh -c '/bin/echo 0 > /sys/devices/pwm-fan/target_pwm'
User=root
StandardOutput=journal+console
Restart=always

[Install]
WantedBy=multi-user.target
