# automagic-fan.service
[Unit]
Description=Automagic fan control

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 -u /usr/local/bin/jetsonfanctl
ExecStopPost=/bin/sh -c '/bin/echo 0 > /sys/devices/pwm-fan/target_pwm'
User=root
StandardOutput=journal+console
Restart=always
SyslogIdentifier=automagic-fan

[Install]
WantedBy=multi-user.target
